# ops/com.redwan.genai-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.redwan.genai-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/redwan/ClaudeProjects/coding-projects/ai-news-feed/scripts/run_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/redwan/ClaudeProjects/coding-projects/ai-news-feed</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/redwan/ClaudeProjects/coding-projects/ai-news-feed/state/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/redwan/ClaudeProjects/coding-projects/ai-news-feed/state/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
